# config/consul-template-files/consul-template-db-upstart.conf
# consul-template
# $ sudo cat /var/log/upstart/consul-template.log
# $ sudo watch -n 1 -d cat /etc/nginx/conf.d/app.conf
description "Consul template process"

start on (local-filesystems and net-device-up IFACE=eth1)
stop on runlevel [!12345]

respawn
respawn limit 10 5


exec consul-template -log-level warn -consul=consul.service.consul:8500 -template "/etc/consul-templates/haproxy.conf:/etc/haproxy/haproxy.cfg:service haproxy restart" > /dev/null 2>&1